# system/systemd/cloudflared-thaleos-guild.service
[Unit]
Description=Cloudflare Tunnel - ThaleOS Guild
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cloudflared tunnel --config /etc/thaleos/cloudflare/config.yml run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
